# awffull.crontab
#
# cron job for the awffull
#
05 * * * *	stats	/usr/sbin/awffull.cron
